# worker/worker.service
[Unit]
Description=Analyze Video Worker
After=network-online.target mnt-nas.automount
Wants=network-online.target

[Service]
Type=simple
User=ghost
WorkingDirectory=/home/ghost/analyze-video
ExecStart=/home/ghost/analyze-video/.venv/bin/uvicorn worker.server:app --host 0.0.0.0 --port 8741 --timeout-keep-alive 120
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
